# build/root/etc/supervisor/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

[program:buildserver]
autostart=true
autorestart=true
command=/app/app
directory=/app
env=PORT="8080",GIN_MODE="release"
user=web
group=web
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
